# virtual-audio-bridge/rootfs/etc/services.d/virtual-audio-bridge/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup Virtual Audio Bridge service
# ==============================================================================

bashio::log.info "Virtual Audio Bridge service stopped"

# Kill any remaining processes
pkill -f "virtual-audio-bridge" || true
pkill -f "parec" || true

# Clean up temporary files
rm -f /tmp/audio-streams/* || true

bashio::log.info "Cleanup completed"